# routers/example/files/etc/crontabs/root
0 4 * * * /etc/scripts/update-ipsets.sh
